# rc.local
#!/bin/bash

zLOGFILE="/Users/xinyi/Desktop/report.txt"
[ ! -f "$zLOGFILE" ] && touch "$zLOGFILE"

zKEXTS="/System/Library/Extensions/IOSCSIArchitectureModelFamily.kext/Contents/PlugIns/IOSCSIMultimediaCommandsDevice.kext"

zTIMESTAMP=`date +"%F %T"`
zMSG="kextunload success"
/sbin/kextunload "$zKEXTS"
if [ $? -ne 0 ]
then
	zMSG="kextunload failed"
fi

echo "[$zTIMESTAMP] $zMSG" >> "$zLOGFILE"
unset zLOGFILE
unset zKEXTS
unset zTIMESTAMP
unset zMSG
